# docker/supervisor.d/horizon.conf
[program:horizon]
process_name=%(program_name)s
command=php /usr/share/nginx/api/artisan horizon
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/usr/share/nginx/api/horizon.log
stopwaitsecs=3600